# devops/mac/ai.openclaw.backup-verify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.backup-verify</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
export PATH="/opt/homebrew/bin:$PATH"
export RESTIC_PASSWORD_FILE="$HOME/.openclaw/restic-password"
REPO="$HOME/openclaw-backups"
LOG="$HOME/.openclaw/logs/backup-verify.log"

echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) [$(hostname -s)] Backup verification starting" >> "$LOG"

if restic -r "$REPO" check --read-data-subset=10% 2>&amp;1 >> "$LOG"; then
  SNAPSHOTS=$(restic -r "$REPO" snapshots --latest 1 --compact 2>&amp;1 | tail -3 | head -1)
  echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) [$(hostname -s)] Verification passed. Latest: $SNAPSHOTS" >> "$LOG"
else
  echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) [$(hostname -s)] VERIFICATION FAILED" >> "$LOG"
  ADMIN_FILE="$HOME/.openclaw/health-check-admin"
  if [ -f "$ADMIN_FILE" ]; then
    ADMIN_CMD=$(sed -n '2p' "$ADMIN_FILE")
    MSG="[$(hostname -s)] Backup verification FAILED. Check ~/.openclaw/logs/backup-verify.log"
    eval "$(echo "$ADMIN_CMD" | sed "s/{MESSAGE}/$MSG/")"
  fi
fi
        </string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/openclaw/backup-verify.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/backup-verify.log</string>
</dict>
</plist>
